# crontab
@reboot bash ~/stereorpicapture/runcamera.sh; sudo shutdown now
